# systemd/snasync.timer
[Unit]
Description=Snasync backup nightly

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
